# dpkg/prerm
#!/bin/sh

systemctl stop p1monitor.service
systemctl disable p1monitor.service
